# init.d/steambox
#!/sbin/runscript

depend() {
	after mtab sabayonlive
	before xdm
	need dbus
	need net
}

start() {
	ebegin "Configuring Sabayon Steam Box if needed"
	/usr/libexec/sabayon-steambox.sh
	eend $?
}

