# init.d/steambox
#!/sbin/runscript

command="/usr/libexec/sabayon-steambox.sh"
description="Sabayon SteamBox service"

depend() {
	after mtab sabayonlive
	before xdm
	need dbus
	need net
}

